# crontab
34 12 * * * python /usr/src/django-scraper/app/db_dump.py
